# roles/repmgr_cluster/templates/postgresql.service.j2
# PostgreSQL systemd service unit file
# Generated by Ansible for cluster {{ cluster_name }}

[Unit]
Description=PostgreSQL {{ postgresql_version }} database server ({{ cluster_name }})
Documentation=https://www.postgresql.org/docs/{{ postgresql_version }}/static/
After=syslog.target
After=network.target

[Service]
Type=notify

User={{ postgresql_user }}
Group={{ postgresql_group }}

# Location of database directory
Environment=PGDATA={{ postgresql_data_dir }}
Environment=PGPORT={{ postgresql_port }}

# Where to send early-startup messages from the server (before the logging
# options of postgresql.conf take effect)
# This is normally controlled by the global default set by systemd
# StandardOutput=syslog

# Disable OOM kill on the postmaster
OOMScoreAdjust=-1000

ExecStart=/usr/pgsql-{{ postgresql_version }}/bin/postgres -D {{ postgresql_data_dir }} -p {{ postgresql_port }}
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT

# Do not set any timeout value, so that systemd will not kill postmaster
# during crash recovery.
TimeoutSec=0

# Set resource limits
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target