# roles/repmgr_cluster/templates/repmgrd.service.j2
# repmgrd systemd service unit file
# Generated by Ansible for cluster {{ cluster_name }}

[Unit]
Description=repmgr daemon for PostgreSQL {{ postgresql_version }} ({{ cluster_name }})
Documentation=https://repmgr.org/
After={{ postgresql_service_name }}.service
Requires={{ postgresql_service_name }}.service

[Service]
Type=forking
User={{ postgresql_user }}
Group={{ postgresql_group }}

# Environment
Environment=PGUSER={{ postgresql_user }}
Environment=PGPORT={{ postgresql_port }}
Environment=PGDATA={{ postgresql_data_dir }}

# Command and configuration
ExecStart=/usr/pgsql-{{ postgresql_version }}/bin/repmgrd -d -f /etc/repmgr/{{ postgresql_version }}/repmgr.conf
ExecReload=/bin/kill -HUP $MAINPID

# Process management
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=30

# Restart policy
Restart=always
RestartSec=10

# Create PID file directory
RuntimeDirectory=repmgr-{{ postgresql_version }}
PIDFile=/run/repmgr-{{ postgresql_version }}/repmgrd.pid

[Install]
WantedBy=multi-user.target